# data/vnc.socket
[Unit]
Before=vnc.service

[Socket]
ListenStream=5900
BindIPv6Only=both
Accept=no
